# modules/kickstandproject/files/graphite/etc/init/statsd.conf
# zuul
#
# Upstart control file for zuul
#

description "Zuul"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
    if [ ! -d /var/run/statsd ]; then
        mkdir -p /var/run/statsd
        chown statsd: /var/run/statsd
    fi
end script

post-stop script
    rm -rf /var/run/statsd
end script

script
    exec start-stop-daemon --start --pidfile /var/run/statsd/statsd.pid -c statsd --exec /usr/bin/node /opt/statsd/stats.js /etc/statsd/config.js --
end script
