# modules/kickstandproject/files/base/etc/cron.d/swapon
# This file is generated by Puppet
#  DO NOT EDIT THIS FILE BY HAND -- YOUR CHANGES WILL BE OVERWRITTEN
#
@reboot root /sbin/swapon /mnt/swap.1
